# rootfs/etc/services.d/model/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Model
# ==============================================================================

TOKEN=$(bashio::config 'token')

bashio::log.info "Starting Model..."

exec python3 /usr/bin/model.py \
    "${TOKEN}" 