# rootfs/etc/services.d/model/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

## Run your program
TOKEN=$(bashio::config 'token')
ENTITIES=$(bashio::config 'entities')
OUT_ENTITIES=$(bashio::config 'out_entities')

bashio::log.info "Starting Model..."

exec python3 /usr/bin/pre_model.py \
    "${TOKEN}" "${ENTITIES}" "${OUT_ENTITIES}"
